# doc/nats-server/etc_systemd_system_nats.service
[Unit]
Description=NATS Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nats-server -c /etc/nats/nats.conf
Restart=on-failure
User=nats
Group=nats
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
